# UNITS/group-template.target
[Unit]
Description=VMMGR %%GROUP%% mgmt
RequiredBy=VMMgr_.target
