# UNITS-Templates/VMMgr_.target
[Unit]
Description=VMMGR VM management
Requires=multi-users.target
